# autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Exec=/home/rdr/.config/polybar/launch.sh
X-GNOME-Autostart-enabled=true
